# fake-c8y/s6-rc.d/mosquitto/run
#!/command/execlineb -P
mosquitto -c /mosquitto/config/mosquitto.conf