# tedge/s6-rc.d/tedge-mapper-c8y/run
#!/command/with-contenv sh
tedge-mapper c8y